# postgres_to_es/etl_pipeline/etc/etl.cron
*/1 * * * * python manager.py  --config etc/config.ini --entity person
*/1 * * * * python manager.py  --config etc/config.ini --entity film_work
*/1 * * * * python manager.py  --config etc/config.ini --entity genre
